# root/etc/s6/strider/finish
#!/bin/bash
echo "strider execution finished that shouldnt happen lets check for errors runned"
