# root/etc/s6/strider/finish
#!/bin/bash
echo "haproxy finish script runned"
