# root/etc/s6/.s6-svscan/finish
#!/bin/bash
echo "Exited and runned finish"
